# var/service/sophie/run
#!/bin/sh
# $Header:  Exp $

# app-antivirus/sophie supervise script 

exec /opt/bin/sophie -d &>/dev/null

