# var/service/xfs/run
#!/bin/sh
# $Header:  Exp $

# x font server supervise script

exec /usr/X11R6/bin/xfs -droppriv

