# var/service/laboule-tarpit/run
#!/bin/sh

exec /local/adm/laboule/bin/laboule-tarpit --server

