# var/service/sshd/run
#!/bin/sh
# $Header:  Exp $

# net-misc/openssh supervise script

exec /usr/sbin/sshd -D

