# var/service/slim/run
#!/bin/sh
# $Header:  Exp $

exec /usr/bin/slim -nodaemon

