# var/service/junkbuster/run
#!/bin/sh
# $Header:  Exp $ 

exec /usr/sbin/junkbuster /etc/junkbuster/config

