# var/service/portmap/run
#!/bin/bash

exec /sbin/portmap -f -t /var/empty -u 113 -g 111


