# var/service/acpid/run
#!/bin/sh

exec /usr/sbin/acpid --foreground &>/dev/null

