# var/service/named/run
#!/bin/sh
# $Header:  Exp $

# net-dns/bind supervise script

exec /usr/sbin/named -f

