# var/service/opendkim/run
#!/bin/sh
# $Header:  Exp $

exec /usr/sbin/opendkim -f -x /etc/opendkim/opendkim.conf

