# var/service/mysqld/run
#!/bin/sh
# $Header:  Exp $

# dev-db/mysql supervise script

mkdir /var/run/mysqld
chown mysql:mysql /var/run/mysqld

exec /usr/sbin/mysqld --user=mysql --binlog-do-db=mysql

