# var/service/nginx/run
#!/bin/sh

rm -f /run/nginx.pid
#mkdir -p /var/tmp/nginx
exec /usr/sbin/nginx

