# var/service/trophie/run
#!/bin/sh
# $Header:  Exp $

# app-antivirus/trophie supervise script

exec /opt/bin/trophie -d &>/dev/null

