# var/service/tinydns/run
#!/bin/sh
# $Header:  Exp $

# net-dns/djbdns supervise script

exec 2>&1
exec envuidgid tinydns envdir /etc/tinydns/env softlimit -d300000 /usr/bin/tinydns

