# var/service/laboule/run
#!/bin/sh

exec /local/adm/laboule/bin/laboule --server

