# var/service/ekeyd/run
#!/bin/bash

exec /usr/libexec/ekeyd -f /etc/entropykey/ekeyd.conf

