# var/service/smartd/run
#!/bin/bash


exec /usr/sbin/smartd --no-fork

