# var/service/ravmd/run
#!/bin/sh
# $Header:  Exp $

exec /opt/rav/bin/ravmd --user qmailq --group qmail --foreground
