# var/service/ntpd/run
#!/bin/sh

# net-misc/openntpd supervise script

sleep 600

rm -f /run/ntpd.sock
exec /usr/sbin/ntpd -d -s

