# var/service/clockspeed-speed/run
#!/bin/sh
# generated by clockspeed-conf-0.4.5_0
# $Header:  Exp $

# net-misc/clockspeed supervise script

exec 2>&1
envdir ./env /bin/sh -c '/usr/bin/sntpclock ${SERVER_IP} | /usr/bin/clockadd'
exec /usr/bin/clockspeed

