# var/service/icecast/run
#!/bin/sh
# $Header:  Exp $ 

exec su - icecast -c "icecast -c /etc/icecast/icecast.conf"
